# launchd/com.d11nn.codex-discord-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.d11nn.codex-discord-notifier</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USER/bin/codex-discord-notifier</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/codex-discord-notifier.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/codex-discord-notifier.err.log</string>
</dict>
</plist>
